# Whitepapers/firmware/fw_images/XOR_example/_dec.rom.extracted/squashfs-root/etc/rc.local
#!/bin/sh

if [ -e "/configData/rc.post" ]; then
	./configData/rc.post
fi
